# recipes-core/images/script-install/20-sd_mount.rules
KERNEL=="mmcblk0p1" SUBSYSTEM=="block" RUN+="/bin/mount /mnt/sd"
